# systemd/nwafuLogin.timer
[Unit]
Description=NWAFU Network Login Checker

[Timer]
OnActiveSec=5min

Unit=nwafuLogin.service

[Install]
WantedBy=timers.target
